# backup-server.service
[Unit]
Description=Qadhya Backup API Server
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/qadhya
ExecStart=/usr/bin/python3 /opt/qadhya/backup-server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/backup-server.log
StandardError=append:/var/log/backup-server.log

# Sécurité
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
